# dotfiles/local/profile.d/105.linuxify.sh
 if [ "$(uname)" == "Darwin" ]; then
     source ~/.linuxify
 fi
