# dotfiles/local/profile.d/804.fzf-obc.sh
source $HOME/.local/opt/fzf-obc/bin/fzf-obc.bash
