# dotfiles/local/profile.d/899.path.sh
export PATH="$HOME/Desktop/todo:$PATH" 
