# dotfiles/local/profile.d/899.path.sh
export PATH="$PATH:$HOME/bin"
# add my own apps to path
export PATH="$HOME/Desktop/todo:$PATH" 

