# dotfiles/local/profile.d/900.pyenv.sh
export PATH="/home/g/.pyenv/bin:$PATH"
eval "$(pyenv init -)"
eval "$(pyenv virtualenv-init -)"
