# dotfiles/local/profile.d/900.dev_cloud_access_token.sh
alias token_dev="curl --fail -v -d 'username=g&password=bliss&org=backed-analytics' 'https://dev-login.relayr.io/oauth/token?client_id=api-client' 2>&1"
# https://github.com/relayr/developer-resources/blob/master/backend/access_tokens.md
