# dotfiles/local/profile.d/901.relayr_env_vars.sh
export CONSUL_HTTP_TOKEN=fd8d4314-6f68-4f4b-9e39-424bfbe897e5
export RELAYR_CONFIG_ENCRYPTION_KEY=KJ0EFvJPnkxP3Smf6QP2aqG23UVJzC4Kmcx1mZ0D4Ks
